# media-sound/lastmp/files/lastmp
#!/sbin/runscript

depend() {
	need lfmsubmitd
}

start () {
	ebegin "Starting MPD last.fm client: lastmp"
	start-stop-daemon --start -c lastmp -p /var/run/lastmp/lastmp.pid --exec /usr/bin/lastmp
	eend $?
}

stop () {
	ebegin "Stopping MPD last.fm client: lastmp"
	start-stop-daemon --stop --quiet -p /var/run/lastmp/lastmp.pid || eerror "lastmp is not running"
	eend $?
}
