# media-sound/lastmp/files/lfmsubmitd
#!/sbin/runscript

depend() {
	need localmount
	need net
}

start () {
	ebegin "Starting MPD last.fm client: lastfmsubmitd"
	start-stop-daemon --start -c lastmp -p /var/run/lastmp/lastfmsubmitd.pid --exec /usr/bin/lastfmsubmitd
	eend $?
}

stop () {
	ebegin "Stopping MPD last.fm client: lastfmsubmitd"
	start-stop-daemon --stop --quiet -p /var/run/lastmp/lastfmsubmitd.pid
	eend $?
}
