# media-sound/mpd-inotify/files/init_d
#!/sbin/runscript

depend() {
	need mpd
}

checkconfig() {
	if ! [ -f /etc/mpd.conf ]; then
		eerror "Configuration file /etc/mpd.conf does not exist."
		return 1
	fi
    # ugh
    MDIR=/$(grep music_directory /etc/mpd.conf | sed "s:['\"]::g" | cut -d/ -f2-)
    # /ugh
	return 0
}

start () {
	checkconfig || return 1
	ebegin "Starting MPD inotify client"
	start-stop-daemon --start -b --exec /usr/bin/mpd-inotify -- ${MDIR}
	eend $?
}

stop () {
	ebegin "Stopping MPD inotify client"
	start-stop-daemon --stop --exec /usr/bin/mpd-inotify
	eend $?
}
